# provisioning/configs/supervisor/mariadb.conf
[program:mysql]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
stderr_logfile=/var/log/mysql.err.log
stdout_logfile=/var/log/mysql.out.log
